# docker/rootfs/etc/services.d/go2rtc/finish
#!/command/with-contenv bashio
# ==============================================================================
# Take down the S6 supervision tree when go2rtc fails
# ==============================================================================
# shellcheck shell=bash

if [[ "${1}" -ne 0 ]] && [[ "${1}" -ne 256 ]]; then
    bashio::log.warning "go2rtc crashed, halting container..."
    /run/s6/basedir/bin/halt
else
    bashio::log.info "go2rtc stopped, restarting..."
fi
