# docker/rootfs/etc/services.d/go2rtc/run
#!/usr/bin/execlineb -P
/usr/local/go2rtc/sbin/go2rtc -config=./go2rtc.yaml