# docker/rootfs/etc/services.d/go2rtc/run
#!/usr/bin/execlineb -P
if [ -f "/config/go2rtc.yaml" ]; then
    /usr/local/go2rtc/sbin/go2rtc -config=/config/go2rtc.yaml
else
    /usr/local/go2rtc/sbin/go2rtc -config=/usr/local/go2rtc/sbin/go2rtc.yaml
fi