# docker/rootfs/etc/services.d/frigate/finish
#!/command/with-contenv bashio
# ==============================================================================
# Take down the S6 supervision tree when Frigate fails
# ==============================================================================
# shellcheck shell=bash

if [[ "${1}" -ne 0 ]] && [[ "${1}" -ne 256 ]]; then
    bashio::log.warning "Frigate crashed, halting container..."
    /run/s6/basedir/bin/halt
fi

bashio::log.info "Frigate stopped, restarting..."
