# docker/rootfs/etc/services.d/frigate/run
#!/command/with-contenv bash
# shellcheck shell=bash

set -euo pipefail

cd /opt/frigate

exec 2>&1
exec python3 -u -m frigate
