# docker/rootfs/etc/services.d/frigate/run
#!/command/with-contenv bashio
# ==============================================================================
# Runs Frigate
# ==============================================================================
# shellcheck shell=bash

bashio::log.info 'Starting Frigate...'

exec python3 -u -m frigate
