# docker/rootfs/etc/services.d/nginx/finish
#!/command/with-contenv bashio
# ==============================================================================
# Take down the S6 supervision tree when NGINX fails
# ==============================================================================
# shellcheck shell=bash

if [[ "${1}" -ne 0 ]] && [[ "${1}" -ne 256 ]]; then
    bashio::log.warning "NGINX crashed, halting container..."
    /run/s6/basedir/bin/halt
else
    bashio::log.info "NGINX stopped, restarting..."
fi
