# docker/rootfs/etc/services.d/nginx/run
#!/usr/bin/execlineb -P
/usr/local/nginx/sbin/nginx