# docker/rootfs/etc/services.d/nginx/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

exec nginx
