# docker/rocm/rootfs/etc/s6-overlay/s6-rc.d/compile-rocm-models/up
/etc/s6-overlay/s6-rc.d/compile-rocm-models/run
