# docker/rootfs/etc/services.d/jsmpeg/run
#!/usr/bin/execlineb -P
/usr/local/bin/node /opt/frigate/jsmpeg/relay.js 8081 8082